# docker/basic-env/supervisor/conf.d/code-server.conf
[program:code-server]
command=%(ENV_VSCODE_BINARY)s serve-web --accept-server-license-terms --without-connection-token
autorestart=unexpected
autostart=false
priority=1
