# docker/supervisor/conf.d/code-server.conf
[program:code-server]
command=code-server --auth none
autorestart=unexpected
autostart=false
priority=1
